# .config/autostart/caps_escape.desktop
[Desktop Entry]
Type=Application
Name=Keyboard setup
Exec=setxkbmap -option caps:escape
Comment=Configure keys
