# .config/autostart/flameshot.desktop
[Desktop Entry]
Type=Application
Name=Flameshot
Exec=flameshot
Comment=Launch Flameshot
