# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Name=Keyboard setup
Exec=xbindkeys
Comment=Configure keys
